# deploy/dr-indexer.service
[Unit]
Description=DigiRails Protocol Indexer
After=network.target digibyted.service
Wants=digibyted.service

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/path/to/digirails-indexer
Environment=PYTHONPATH=/path/to/digirails-indexer
EnvironmentFile=/path/to/digirails-indexer/.env
ExecStart=/path/to/digirails-indexer/venv/bin/python -m dr_indexer
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
